# templates/trackpoint.timer
[Unit]
Description=Activates trackpoint after boot

[Timer]
OnActiveSec=10
Persistent=true

[Install]
WantedBy=timers.target
